# SDK_AutoTest/SE5901A_SDK_Test_Tool/wireless/rc.local
#!/bin/sh
webs &

. /jffs2/wireless/rc.init
